# buildroot-2021.02.2/board/pine64/pine64/fs/etc/systemd/system/wayland.socket
[Unit]
Description=Wayland Socket
Before=sockets.target

[Socket]
ListenStream=/tmp/wayland-0

