# buildroot-2021.02.2/board/pine64/pine64/fs/etc/systemd/system/wayland.service
[Unit]
Description=weston
Requires=wayland.socket 
After=wayland.socket 

[Service]
EnvironmentFile=/etc/sysconfig/pineroot
ExecStart=/usr/bin/weston -i0 --tty=3 
Restart=always

