# buildroot-2021.02.2/board/pine64/pine64/fs/etc/profile.d/prompt.sh
export PS1='\u@\h [`cat /sys/class/power_supply/axp20x-battery/voltage_now`]: \w \$ '
export XDG_RUNTIME_DIR=/tmp
